# DEBIAN/postinst
#!/bin/bash

echo -e "\n$0\n"

cd /

if [ ! -d ./usr/local/u/bin ];
then
   echo "$0: Abort.  " `realpath ./usr/local/u/bin`" does not exist"
   exit 1
else
   chown -R jjchkconn:jjchkconn ./usr/local/u/bin/jjchkconn*
   chmod -R  550 ./usr/local/u/bin/jjchkconn*
fi

chown jjchkconn:jjchkconn ./home/jjchkconn/.config/uSysIntChkd/*
chmod u+w ./home/jjchkconn/.config/uSysIntChkd/

systemctl daemon-reload

systemctl restart apparmor

# Append to /etc/aliases
cat << EOF >> /etc/aliases
jjchkconn: /home/jjchkconn/alerts/inbound-alerts
alerts: /home/jjchkconn/alerts/inbound-alerts 
EOF

# Update aliases database
newaliases

systemctl enable jjchkconn-uSysIntChkd@jjchkconn.service  || true
systemctl enable jjchkconn-alerts-send@jjchkconn.path || true

nohup sh -c 'sleep 2; systemctl start jjchkconn-uSysIntChkd@jjchkconn.service' &
nohup sh -c 'sleep 2; systemctl start  jjchkconn-alerts-send@jjchkconn.path' &

systemctl status -a "jj*" |cat  || true

echo -e "\n\nREMEMBER TO CUSTOMIZE YOUR SLACK WEBHOOKS PARAMETER IN THE .INI FILE."
echo -e "\nIT IS NOT CONFIGURED RIGHT NOW.\n"
echo -e "\n\nALSO, VERIFY THE .NEXT AND .STEPS FILES."
echo -e "\nInstallation is complete.\n"
