# etc/systemd/system/jjchkconn-alerts-send@.path
[Unit]
Description="Monitor the /home/%i/alerts/inbound-alerts file for changes"

[Path]
PathModified=/home/%i/alerts/inbound-alerts
Unit=jjchkconn-alerts-send@%i.service

[Install]
WantedBy=default.target
