# etc/systemd/system/jjchkconn-alerts-send@.path
[Unit]
Description="Monitor the /home/jjchkconn/alerts/inbound-alerts file for changes"

[Path]
User=%i
Group=%i
PathModified=/home/jjchkconn/alerts/inbound-alerts
Unit=jjchkconn-alerts-send@%i.service

[Install]
WantedBy=default.target
